# src/build/modules/redis/files/etc/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/local/bin/redis-server /redis/conf.d/redis.conf
autostart=true
autorestart=true
