# src/build/modules/redis/files/etc/supervisor/conf.d/redis.conf
[program:redis]
command = /usr/local/bin/redis-server /root/redis.conf
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
