# build/etc/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/local/bin/redis-server /redis-2.8.14/conf.d/redis.conf
autostart=true
autorestart=true